# php/docker/templates/supervisor.d/blackfire.conf
{{ if .Env.BLACKFIRE_SERVER_ID }}
[program:blackfire-agent]
command = blackfire-agent --log-level=3 --config=/usr/local/docker/etc/blackfire-agent.conf
process_name=%(program_name)s
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
{{ end }}