# root/etc/cron.d/wp-cron
* * * * * www-data /usr/loca/bin/wp cron event run --due-now
